# tdx-ml-workload/build/linux/ita-demo.service
[Unit]
Description=Intel Trust Authority Demo Application

[Install]
WantedBy=multi-user.target

[Service]
Type=simple
WorkingDirectory=/opt/ita-demo
ExecStart=/usr/local/bin/ita-demo
ExecReload=/bin/kill -s HUP $MAINPID
Restart=on-failure
TimeoutStartSec=0
PermissionsStartOnly=true
RuntimeDirectory=ita-demo
RuntimeDirectoryMode=0775
