# tdx-ml-workload/build/linux/trustauthority-demo.service
[Unit]
Description=Intel Trust Authority Demo Application

[Install]
WantedBy=multi-user.target

[Service]
Type=simple
WorkingDirectory=/opt/trustauthority-demo
ExecStart=/usr/local/bin/trustauthority-demo
ExecReload=/bin/kill -s HUP $MAINPID
Restart=on-failure
TimeoutStartSec=0
PermissionsStartOnly=true
RuntimeDirectory=trustauthority-demo
RuntimeDirectoryMode=0775
